# docker-nginx/crontab
* *     * * *   root    /etc/my_init.d/z_restart_nginx
