# docker-seafilepro/scripts/seafile.crontab
#0 0 * * * /etc/my_init.d/02_backup.sh
0 0,12 * * * root /usr/local/bin/s3ql_backup.sh >/dev/null 2>&1
